# host_admin/moonraker-admin.service
[Unit]
Description=Moonraker Host Admin Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/moonraker-admin
EnvironmentFile=/opt/moonraker-admin/.env
ExecStart=/usr/bin/python3 -m uvicorn admin_service:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
